# webserver/FAB/files/supervisor/opentrain.conf
[program:opentrain]
user = %(USER)s
directory = %(DJANGO_BASE_DIR)s
command = %(HOME)s/bin/run_gunicorn.sh
stdout_logfile = %(HOME)s/log/run.log
redirect_stderr = true








